# deploy/systemd/lan-game-hub.target
[Unit]
Description=LAN Game Hub (tutti i server di gioco)
Wants=doom-relay.service quake3-server.service

[Install]
WantedBy=multi-user.target
